# systemd/telegram-channels-to-md.service.in
[Unit]
Description=Export Telegram channels to Markdown (python_telegram_channels_to_md)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_DIR@
# Run as specific user (filled by Makefile during install)
User=@RUN_AS_USER@
Environment=HOME=@RUN_AS_HOME@
ExecStart=/usr/bin/env bash -lc 'bash "@REPO_DIR@/scripts/export_and_push.sh" "@REPO_DIR@"'
Environment=PYTHONUNBUFFERED=1
# Optional env file with git/push settings (ignored if absent)
EnvironmentFile=-@REPO_DIR@/systemd/telegram-channels-to-md.env

[Install]
WantedBy=default.target