# systemd/telegram-channels-to-md.service.in
[Unit]
Description=Export Telegram channels to Markdown (python_telegram_channels_to_md)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_DIR@
ExecStart=/usr/bin/env bash -lc 'bash "@REPO_DIR@/scripts/export_and_push.sh" "@REPO_DIR@"'
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target